# web_app/webapp.service
[Unit]
Description=Miles DevOps Case Study — Gunicorn
After=network.target postgresql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/app/miles-studycase/web_app
EnvironmentFile=/app/.env
ExecStartPre=/app/miles-studycase/web_app/venv/bin/miles-challenge-initdb
ExecStart=/app/miles-studycase/web_app/venv/bin/gunicorn -c /app/gunicorn.conf.py miles_challenge.app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target